# root/etc/supervisor.d/supervisord.ini
[supervisord]
nodaemon                = true
logfile                 = /dev/null
logfile_maxbytes        = 0
logfile_backups         = 0
user                    = root
pidfile                 = /run/supervisord.pid

[program:apache2]
command                 = /usr/sbin/httpd -DFOREGROUND
autorestart             = true
killasgroup             = true
stopasgroup             = true
redirect_stderr         = true
stdout_logfile          = /proc/self/fd/1
stdout_logfile_maxbytes = 0
stdout_logfile_backups  = 0

[program:backuppc]
command                 = /usr/local/share/backuppc/bin/BackupPC
autorestart             = true
user                    = apache
killasgroup             = true
stopasgroup             = true
redirect_stderr         = true
stdout_logfile          = /proc/self/fd/1
stdout_logfile_maxbytes = 0
stdout_logfile_backups  = 0